# upcase-podman@.socket
[Unit]
Description=upcase service socket for %I

[Socket]
#ListenStream=%h/upcase-podman-sock.%i

[Install]
WantedBy=default.target
